# bootstrap/roles/jaeger/templates/jaeger.service.j2
[Unit]
Description=Jaeger
After=docker.service
Requires=docker.service

[Service]
Restart=always
RestartSec=5
ExecStartPre=-/usr/bin/docker stop jaeger
ExecStartPre=-/usr/bin/docker rm jaeger
ExecStart=/usr/bin/docker run \
  --name jaeger \
  -p {{ jaeger_port_collector }}:14268 \
  -p {{ jaeger_port_ui }}:16686 \
  --cpus={{ jaeger_cpus }} \
  --memory={{ jaeger_memory }} \
  -e TZ={{ jaeger_tz }} \
  -v {{ jaeger_config_dir }}/config.yaml:/etc/jaeger/config.yaml:ro \
  {{ jaeger_image }} \
  --config=file:/etc/jaeger/config.yaml
ExecStop=/usr/bin/docker stop jaeger

[Install]
WantedBy=multi-user.target
