# bootstrap/roles/l8opensim/templates/l8opensim.service.j2
[Unit]
Description=L8 OpenSim SNMP simulator
After=docker.service
Requires=docker.service

[Service]
Restart=always
RestartSec=5
ExecStartPre=-/usr/bin/docker stop l8opensim
ExecStartPre=-/usr/bin/docker rm l8opensim
ExecStart=/usr/bin/docker run \
  --name l8opensim \
  --privileged \
  --network host \
  --device /dev/net/tun:/dev/net/tun \
  {{ l8opensim_image }} \
  -auto-start-ip {{ l8opensim_auto_start_ip }} \
  -auto-count {{ l8opensim_auto_count }} \
  -auto-netmask {{ l8opensim_auto_netmask }} \
  -port {{ l8opensim_web_port }}
ExecStop=/usr/bin/docker stop l8opensim

[Install]
WantedBy=multi-user.target
